# chef-solo/cookbooks/h2/templates/default/init-h2.sh.erb
#!/bin/bash
# description: Tomcat Start Stop Restart
# processname: tomcat
# chkconfig: 234 20 80
PATH=/bin:/usr/bin:/sbin:/usr/sbin:/usr/local/tomcat/bin
# The first existing directory is used for JAVA_HOME (if JAVA_HOME is not defined in $DEFAULT)
JDK_DIRS="/usr/lib/jvm/java-6-openjdk /usr/lib/jvm/java-6-sun /usr/lib/jvm/java-1.5.0-sun /usr/lib/j2sdk1.5-sun /usr/lib/j2sdk1.5-ibm"

# Look for the right JVM to use
for jdir in $JDK_DIRS; do
	if [ -r "$jdir/bin/java" -a -z "${JAVA_HOME}" ]; then
	JAVA_HOME="$jdir"
	fi
done
export JAVA_HOME

PATH=$JAVA_HOME/bin:$PATH
export PATH  

base=/usr/local/h2
classpath=.:/usr/local/h2/h2.jar
netstat=/bin/netstat
ps="ps axwww"
server=org.h2.tools.Server
webport=8082
defaultports="(5435\|${webport}\|9092\|9094)"

getpid() {
  pid=`${ps} | grep java | grep ${server} | awk '{ print $1 };'`
}

start() {
  java -cp ${classpath} ${server} -tcp -web -baseDir ${base} &
  READY=0
  while [ $READY -eq 0 ]
    do
      sleep 1
      READY=`${netstat} -an | grep ${webport} | wc -l`
    done
}

stop() {
  getpid
  if [ "${pid}x" = "x" ]
    then echo "No H2 Server running."
  else
    echo "Killing H2 Server process id: ${pid}"
    kill -9 ${pid}
  fi
}

status() {
  getpid
  if [ "${pid}x" = "x" ]
    then 
     echo "No H2 Server running."
     exit 1
  else
    echo "H2 Server has process id: ${pid}"
    ${netstat} -an | grep ${defaultports} | grep -i listen
  fi
}

restart() {
    stop
    start
}

help() {
  java -cp ${classpath} ${server} -?
}

case "$1" in
  start)
    start
    ;;
  stop)
    stop
    ;;
  restart)
    restart
    ;;
  status)
    status
    ;;
  help)
    help
    ;;
  *)
    echo $"Usage: $0 {start|stop|restart|status|help}"
    exit 1
esac

exit $?
